# misc/tmpfilesd.cron
# execute tmpfilesd cleanup after boot and then daily

@boot	/sbin/tmpfilesd --clean
@daily	/sbin/tmpfilesd --clean
